# drive/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=-/sbin/agetty --skip-login --noclear --noissue --login-options "-f root" %I $TERM
